# .Docker/root/etc/service/apache/run
#!/bin/sh

set -x

exec 2>&1

rm -r /var/run/apache2/httpd.pid
exec /usr/sbin/httpd -D FOREGROUND